# upstart/nera_bot_ros.service
# THIS IS A GENERATED FILE, NOT RECOMMENDED TO EDIT.

[Unit]
Description="bringup nera_bot_ros"
After=network.target
Wants=network.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=robot
Group=robot
Environment="ROS_DOMAIN_ID=88"
Environment="ROS_DISTRO=jazzy"
ExecStart=/usr/sbin/nera_bot_ros-start
Restart=on-failure
RestartSec=5
SuccessExitStatus=0 129 200

[Install]
WantedBy=multi-user.target
